# data/_system/supervisord.conf
[unix_http_server]
file=/tmp/museon-supervisor.sock
chmod=0700

[supervisord]
logfile=/Users/ZEALCHOU/MUSEON/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/museon-supervisord.pid
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/museon-supervisor.sock

[program:museon-gateway]
command=/Users/ZEALCHOU/MUSEON/.venv/bin/python -m museon.gateway.server
directory=/Users/ZEALCHOU/MUSEON
environment=MUSEON_HOME="/Users/ZEALCHOU/MUSEON",PYTHONPATH="/Users/ZEALCHOU/MUSEON/src",HOME="/Users/ZEALCHOU",USER="ZEALCHOU",PATH="/usr/sbin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin"
autostart=true
autorestart=unexpected
exitcodes=0
startretries=5
startsecs=10
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/Users/ZEALCHOU/MUSEON/logs/gateway.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/Users/ZEALCHOU/MUSEON/logs/gateway.err
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
